# bin/rikerio@.service
[Unit]
Description=Starting RikerIO Profile %I

[Service]
Type=notify
NotifyAccess=all
Group=rikerio
ExecStart=/usr/bin/rikerio-server --id=%i --size=4096

[Install]
WantedBy=multi-user.target
